# systemd/foxess-soak.timer
[Unit]
Description=Nightly FoxESS soak tests

[Timer]
OnCalendar=*-*-* 12:00:00 UTC
Persistent=true
RandomizedDelaySec=300

[Install]
WantedBy=timers.target
