# systemd/foxess-soak.service
[Unit]
Description=FoxESS Control soak tests
After=network-online.target

[Service]
Type=oneshot
WorkingDirectory=%h/git/fox
ExecStart=/bin/bash -c '\
  set -euo pipefail; \
  REPO="$HOME/git/fox"; \
  TAG=$(git -C "$REPO" describe --tags --abbrev=0); \
  COMMIT=$(git -C "$REPO" rev-parse "$TAG"); \
  SHORT=$(git -C "$REPO" rev-parse --short "$TAG"); \
  BRANCH="$TAG"; \
  TS=$(date +%%Y%%m%%d_%%H%%M%%S); \
  RUN_DIR="$REPO/test-artifacts/soak/runs/${TS}_${SHORT}"; \
  WORKTREE="$REPO/test-artifacts/soak/.worktree-$$"; \
  mkdir -p "$RUN_DIR"; \
  echo "commit: $COMMIT" > "$RUN_DIR/meta.txt"; \
  echo "branch: $BRANCH" >> "$RUN_DIR/meta.txt"; \
  echo "mode: real-time" >> "$RUN_DIR/meta.txt"; \
  echo "started: $(date -Iseconds)" >> "$RUN_DIR/meta.txt"; \
  git -C "$REPO" worktree add --detach "$WORKTREE" "$TAG" 2>&1; \
  trap "git -C \"$REPO\" worktree remove --force \"$WORKTREE\" 2>/dev/null || true" EXIT; \
  cd "$WORKTREE"; \
  export SOAK_ARTIFACT_DIR="$RUN_DIR"; \
  /usr/bin/python3 -m pytest tests/soak/ -m soak \
    --tb=long -v --override-ini="addopts=-n auto" \
    > "$RUN_DIR/pytest.log" 2>&1; \
  EC=$?; \
  echo "exit_code=$EC" >> "$RUN_DIR/pytest.log"; \
  echo "finished: $(date -Iseconds)" >> "$RUN_DIR/meta.txt"; \
  echo "exit_code: $EC" >> "$RUN_DIR/meta.txt"; \
  exit $EC'
TimeoutStartSec=12h
Nice=10
IOSchedulingClass=best-effort
IOSchedulingPriority=7
