# systemd/foxess-soak.service
[Unit]
Description=FoxESS Control soak tests
After=network-online.target

[Service]
Type=oneshot
WorkingDirectory=%h/git/fox
ExecStart=%h/git/fox/systemd/foxess-soak-run.sh
TimeoutStartSec=12h
Nice=10
IOSchedulingClass=best-effort
IOSchedulingPriority=7
